# psql-backup.timer
[Unit]
Description=Postgresql backup

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
